# .config/systemd/user/mconnect.service
[Unit]
Description=KDEConnect Protocol implemented in Vala

[Service]
Type=simple
ExecStart=/usr/bin/mconnect -d

[Install]
WantedBy=default.target
