# config/systemd/user/cbatticon.service
[Unit]
Description=cbatticon

[Service]
ExecStart=/usr/bin/cbatticon
