# config/systemd/user/fcitx5.service
[Unit]
Description=Fcitx 5

[Service]
ExecStart=/usr/bin/fcitx5
