# config/systemd/user/gnome-keyring-secrets.service
[Unit]
Description=Secret Storage Service

[Service]
ExecStart=/usr/bin/gnome-keyring-daemon --start --components=secrets
