# config/systemd/user/gnome-keyring-ssh.service
[Unit]
Description=SSH Key Agent

[Service]
ExecStart=/usr/bin/gnome-keyring-daemon --start --components=ssh
