# config/systemd/user/gnome-keyring-pkcs11.service
[Unit]
Description=Certificate and Key Storage

[Service]
ExecStart=/usr/bin/gnome-keyring-daemon --start --components=pkcs11
