# config/systemd/user/gnome-keyring.service
[Unit]
Description=GNOME Keyring Components

[Service]
ExecStart=/usr/bin/gnome-keyring-daemon --start --foreground
