# config/systemd/user/nitrogen.service
[Unit]
Description=Restore wallpapers with Nitrogen

[Service]
Type=oneshot
ExecStart=/usr/bin/nitrogen --restore
