# config/systemd/user/pasystray.service
[Unit]
Description=PulseAudio System Tray

[Service]
ExecStart=/usr/bin/pasystray
