# config/systemd/user/picom.service
[Unit]
Description=picom

[Service]
ExecStart=/usr/bin/picom
