# config/systemd/user/sxhkd.service
[Unit]
Description=sxhkd

[Service]
ExecStart=/usr/bin/sxhkd
