# config/systemd/user/touchegg.service
[Unit]
Description=Touchégg

[Service]
ExecStart=/usr/bin/touchegg
