# config/systemd/user/touchegg.service
[Unit]
Description=Touchégg
After=graphical.target

[Service]
ExecStart=/usr/bin/touchegg
